# health_eye_celery.ini
[supervisord]
nodaemon=true

[program:health_eye_celery]
command=celery -A manage.celery worker --loglevel=info
directory=/app
user=root
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
stdout_logfile=/app/log/celery.log
stderr_logfile=/app/log/celery.err

