# health_eye_api.ini
[supervisord]
nodaemon=true

[program:health_eye_api]
command=gunicorn -k gevent -w3 -b0.0.0.0:9000 manage:wtf_app
directory=/app
user=root
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
stdout_logfile=/app/log/gunicorn.log
stderr_logfile=/app/log/gunicorn.err

